# system-wide-scripts/profile.d/lm-env.sh
# /etc/profile.d/lm-env.sh
# LM System Environment Variables
# Version: 1.0
# Author: Echocolate from Sakura-Moon Studio

# export PATH="$PATH:/usr/local/lm/bin/"
export lm_sshPort='22'
export lm_wwwUser='www'
export lm_dryRun='0'